# deb/deb-pkg/DEBIAN/postinst
#!/bin/bash

sudo update-rc.d bottlefs defaults
/etc/init.d/bottlefs start
