# deb/deb-pkg/DEBIAN/prerm
#!/bin/bash

/etc/init.d/bottlefs stop
update-rc.d -f bottlefs remove
